# install/lpanel.service
[Unit]
Description=L-Panel FastAPI Service
After=network.target

[Service]
User=root
WorkingDirectory=/opt/l-panel
ExecStart=/usr/local/bin/uvicorn app.main:app --host 0.0.0.0 --port 2096
Restart=always
RestartSec=5
Environment="PATH=/usr/local/bin:/usr/bin:/bin:/usr/sbin:/usr/local/sbin"

[Install]
WantedBy=multi-user.target
